# engines/registration_sharing/package/rmt-server-trim-cache.service
[Unit]
Description=Trim expired cache files generated by RMT
After=rmt-server.service
Wants=rmt-server.service

[Service]
Type=oneshot
Restart=no
Environment="LOG_TO_JOURNALD=1" "LANG=en"
EnvironmentFile=/var/lib/rmt/rmt-cache-trim.sh
ExecStart=find ${REPOSITORY_CLIENT_CACHE_DIRECTORY} -mmin +${REPOSITORY_CACHE_EXPIRY_MINUTES} -type f -delete
ExecStart=find ${REGISTRY_CLIENT_CACHE_DIRECTORY} -mmin +${REGISTRY_CACHE_EXPIRY_MINUTES} -type f -delete

[Install]
WantedBy=multi-user.target
