# engines/registration_sharing/package/rmt-server-trim-cache.service
[Unit]
Description=Trim expired cache files generated by RMT
After=rmt-server.service
Wants=rmt-server.service

[Service]
Type=simple
Restart=no
Environment="LOG_TO_JOURNALD=1" "LANG=en"
ExecStart=find /var/lib/rmt/tmp/cache -mmin +20 -type f -delete

[Install]
WantedBy=multi-user.target
