# package/rmt.target
[Unit]
Description=Repository Mirroring Tool
Requires=mysql.service
Requires=rmt-migration.service
Requires=rmt.service

[Install]
WantedBy=multi-user.target
